# pushman/deploy.bak/pushman_supervisor.ini
[program:pushman]
command=/home/diors/env/bin/uwsgi --ini /home/diors/work/pushman/pushman/deploy/pushman.ini
autostart=true
autorestart=true
startretries=3
exitcodes=0,1,2
stopsignal=KILL
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor_pushman_access.log
stderr_logfile=/var/log/supervisor_pushman_error.log
user=diors